# systemd-units/odasrv.target
[Unit]
Description=Group target for odamex server instances

[Install]
WantedBy=multi-user.target
